# agent-a-infrastructure/systemd/vrops-alert-caller.service
[Unit]
Description=vROps Alert AutoCaller
Documentation=https://github.com/your-org/vrops-alert-caller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vrops-alert
Group=vrops-alert
WorkingDirectory=/opt/vrops-alert-caller
Environment=PATH=/opt/vrops-alert-caller/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1

# workers=1 threads=4：
# 確保去重快取 (in-memory) 與 SipEngine 單例在同一 Process，
# 多 worker 會導致狀態不共享（去重失效、多個 SIP 同時註冊衝突）
ExecStart=/opt/vrops-alert-caller/venv/bin/gunicorn \
    --bind 0.0.0.0:5000 \
    --workers 1 \
    --threads 4 \
    --timeout 120 \
    --access-logfile /opt/vrops-alert-caller/logs/access.log \
    --error-logfile /opt/vrops-alert-caller/logs/gunicorn-error.log \
    webhook_server:app

Restart=always
RestartSec=10

StandardOutput=append:/opt/vrops-alert-caller/logs/service.log
StandardError=append:/opt/vrops-alert-caller/logs/service-error.log

# 資源限制
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
